# rootfs/usr/lib/systemd/system/bootc-fetch-apply-updates.service.d/staged-updates.conf
[Service]
ExecStart=
ExecStart=/usr/bin/bootc upgrade --quiet
