# rootfs/usr/lib/systemd/system/var-mnt-media.mount
[Unit]
Description=Media disk mount

[Mount]
What=MEDIA_DISK_PLACEHOLDER
Where=/var/mnt/media
Type=xfs
Options=defaults
