# rootfs/usr/lib/systemd/system/var-mnt-media.automount
[Unit]
Description=Automount for media disk
Before=local-fs.target

[Automount]
Where=/var/mnt/media

[Install]
WantedBy=local-fs.target
