# rootfs/usr/lib/systemd/system/var-mnt-backup.automount
[Unit]
Description=Automount for backup disk
Before=local-fs.target

[Automount]
Where=/var/mnt/backup
TimeoutIdleSec=5min

[Install]
WantedBy=local-fs.target
