# rootfs/usr/lib/systemd/system/var-mnt-backup.mount
[Unit]
Description=Backup disk mount

[Mount]
What=BACKUP_DISK_PLACEHOLDER
Where=/var/mnt/backup
Type=xfs
Options=defaults
